# supervisor.d/deadd.conf
[program:notifications]
command=deadd-notification-center
autostart=true
autorestart=true
