# supervisor.d/deadd.conf
[program:deadd]
command=deadd-notification-center
autostart=true
autorestart=true
