# supervisor.d/dunst.conf
[program:dunst]
command=dunst
autostart=false
autorestart=true
