# supervisor.d/dunst.conf
[program:notifications]
command=dunst
autostart=true
autorestart=true
environment=DISPLAY=:0
