# supervisor.d/dunst.conf
[program:notifications]
command=dunst
autostart=false
autorestart=true
