# supervisor.d/maestral.conf
[program:dropbox]
command=%(ENV_HOME)s/.local/virtualenvs/maestral/bin/maestral gui
autostart=true
environment=DISPLAY=":0"
