# supervisor.d/picom.conf
[program:compositor]
command=picom --experimental-backends
autostart=true
