# supervisor.d/picom.conf
[program:compositor]
command=picom
autostart=false
